# boot_status.service
[Unit]
Description=Send system status to Discord on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
Environment=HOME=__HOME__
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=__HOME__/scripts/pi_status_scripts/boot_status.sh

[Install]
WantedBy=multi-user.target
